# vola.service
[Unit]
Description=Vola Autonomous Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vola
Group=vola
WorkingDirectory=/home/vola
Environment=VOLA_HOME=/home/vola
ExecStart=/home/vola/daemon/venv/bin/python3 /home/vola/daemon/runner.py
Restart=on-failure
RestartSec=30

StandardOutput=journal
StandardError=journal

ProtectSystem=strict
ReadWritePaths=/home/vola
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
